# fs_monitor_supervisor.conf
[program:redroid-fs-monitor]
directory=/home/redroid/video_saves/redroid-fs-monitor
command=python fs_monitor.py
autostart=true
autorestart=true
